# riin-summaly/docs/deploy-examples/summaly.service.example
# systemd unit example for summaly (TOML config + bin/summaly-server.ts via tsx).
#
# Replace placeholders:
#   /opt/summaly                       — install directory (git checkout で `pnpm install` 済み)
#   /etc/summaly/config.toml           — TOML 設定ファイル（summaly-config.example.toml をコピーして編集）
#   summaly                             — service user / group (create with: useradd -r -s /usr/sbin/nologin summaly)
#
# Install: place at /etc/systemd/system/summaly.service, then:
#   systemctl daemon-reload && systemctl enable --now summaly
#
# `[server]` host / port は config.toml で指定する（fastify-cli 由来の --address / --port は廃止）。

[Unit]
Description=summaly preview server
Documentation=https://github.com/misskey-dev/summaly
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=summaly
Group=summaly
WorkingDirectory=/opt/summaly
ExecStart=/usr/bin/env node_modules/.bin/tsx bin/summaly-server.ts /etc/summaly/config.toml
Restart=on-failure
RestartSec=5

# IPv4 のみ / IPv6 のみで運用する場合は SUMMALY_FAMILY を指定する
# Environment=SUMMALY_FAMILY=4

# プライベート IP 解放はテスト専用。本番では絶対に有効化しない
# Environment=SUMMALY_ALLOW_PRIVATE_IP=true

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
